# debian/spreedbox-mydyns-update.upstart
description "Periodically update mydyns IP"

start on net-device-up
stop on runlevel [!2345]

env config=/etc/default/spreedbox-mydyns-update
env updater=/opt/spreedbox/bin/mydyns-update

pre-start script
[ -e "$config" ] || { stop; exit 0; }
end script

script
while [ 1 ]
do
	. $config
	$updater
	if [ -f "$CACHE/$PREFIX.4" -o -f "$CACHE/$PREFIX.6" ]; then
		sleep 300
	else
		sleep 5
	fi
done
end script
