# debian/spreedbox-firstboot.upstart
# spreedbox-firstboot
#   Runs long running first boot actions.

description "Spreedbox first boot"

start on stopping firstboot
stop on runlevel [!2345]

umask 077

script

	# Create dhparam - this takes a very long time.
	/opt/spreedbox/bin/dhparam-create

end script
